# mcp-server.service
[Unit]
Description=VyOS MCP Server
After=network.target vyos-router.target

[Service]
Type=simple
ExecStart=/config/user-data/vyos-mcp-go
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
