# mcp-server.service
[Unit]
Description=VyOS MCP Server
After=network.target vyos-router.target

[Service]
Type=simple
ExecStart=/config/user-data/vyos-mcp-go
EnvironmentFile=/config/user-data/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
